# scripts/daily/com.josh.wordle-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.josh.wordle-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/josh/code/wordle-bot/scripts/daily/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/josh/code/wordle-bot/scripts/daily/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/josh/code/wordle-bot/scripts/daily/logs/launchd.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
